# services/strongswan/templates/tc-rules.service.j2
[Unit]
Description=Apply Traffic Control Rules

[Service]
Type=oneshot
ExecStart=/usr/local/bin/apply_tc_rules.sh

[Install]
WantedBy=multi-user.target